# com.codex-quota-injector.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.codex-quota-injector.agent</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/zhangzhengwei/code/codex-quota-injector/.runtime/Codex Quota Injector</string>
    <string>/Users/zhangzhengwei/code/codex-quota-injector/src/cli.mjs</string>
    <string>inject</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>CODEX_QUOTA_CDP_PORT</key>
    <string>9229</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>/Users/zhangzhengwei/Library/Logs/Codex Quota Injector/injector.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/zhangzhengwei/Library/Logs/Codex Quota Injector/injector.log</string>
</dict>
</plist>
